# launchd/com.yarbo.chicago-after-sales-dashboard-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yarbo.chicago-after-sales-dashboard-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/warehouse-flow-dashboard/scripts/sync_feishu_sheet.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/warehouse-flow-dashboard</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/warehouse-flow-dashboard/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/warehouse-flow-dashboard/logs/launchd.err.log</string>
</dict>
</plist>
