# launchd/com.yarbo.chicago-dashboard-pages-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yarbo.chicago-dashboard-pages-publish</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/chicago-dashboard-publish/scripts/sync_and_publish.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/logs/github-pages-publish.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mac/Library/Application Support/YarboChicagoDashboard/logs/github-pages-publish.err.log</string>
</dict>
</plist>
